# dist-material/config-examples/skywalking.ini
[program:skywalking]
command = /data/apps/opt/skywalking/bin/oapService.sh start-foreground
directory = /data/apps/opt/skywalking
autostart = true
autorestart = true
startsecs = 10
startretries = 3
stopsignal=TERM
killasgroup = true
user = sw
environment=
    SW_JDBC_URL="jdbc:mysql://10.0.0.2:3307/swtest",
    SW_DATA_SOURCE_USER="swtest",
    SW_DATA_SOURCE_PASSWORD="swtest123",
    OAP_LOG_DIR="/data/apps/log/skywalking",
    SW_STORAGE="mysql",
    JAVA_OPTS="-Xms2048M -Xmx2048M",
    SW_RECEIVER_ZIPKIN="default",
    PATH="/data/apps/opt/mysql/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin",
    redirect_stderr = false
stdout_logfile=/data/apps/log/skywalking/%(program_name)s_std.log
stdout_logfile_maxbytkafka = 300MB
stdout_logfile_backups= 20
stdout_events_enabled = false
stderr_logfile=/data/apps/log/skywalking/%(program_name)s_err.log
stderr_logfile_maxbytkafka=300MB
stderr_logfile_backups=10
stderr_capture_maxbytkafka=3MB
stderr_events_enabled= false


[program:skywalking-webapp]
command = /data/apps/opt/skywalking/bin/webappService.sh start-foreground
directory = /data/apps/opt/skywalking
autostart = true
autorestart = true
startsecs = 10
startretries = 3
stopsignal=TERM
killasgroup = true
user = sw
environment=
    WEBAPP_LOG_DIR="/data/apps/log/skywalking",
    JAVA_OPTS="-Xms2048M -Xmx2048M",
    PATH="/data/apps/opt/mysql/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin",
redirect_stderr = false
stdout_logfile=/data/apps/log/skywalking/%(program_name)s_std.log
stdout_logfile_maxbytkafka = 300MB
stdout_logfile_backups= 20
stdout_events_enabled = false
stderr_logfile=/data/apps/log/skywalking/%(program_name)s_err.log
stderr_logfile_maxbytkafka=300MB
stderr_logfile_backups=10
stderr_capture_maxbytkafka=3MB
stderr_events_enabled= false